# scripts/bot-live.service
# Copy file này vào VPS: sudo cp scripts/bot-live.service /etc/systemd/system/
# Rồi: sudo systemctl daemon-reload && sudo systemctl enable bot-live && sudo systemctl start bot-live
# Đường dẫn bên dưới dùng khi clone về /root/Bot_VPS_Realtime_ver1_11_3_26 (user root)

[Unit]
Description=Bot Trading Paper/Live (Telegram + Web)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/Bot_VPS_Realtime_ver1_11_3_26
ExecStart=/root/Bot_VPS_Realtime_ver1_11_3_26/.venv/bin/python run_live.py
Restart=always
RestartSec=10
EnvironmentFile=/root/Bot_VPS_Realtime_ver1_11_3_26/.env

[Install]
WantedBy=multi-user.target
